# php-stacks/.docker/kibana/config/supervisord.conf
[supervisord]
nodaemon=true

[program:kibana]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/opt/kibana/bin/kibana