# lep-stack/.docker/kibana/config/supervisord.conf
[supervisord]
nodaemon=true

[program:kibana]
command=/opt/kibana/bin/kibana